# com.user.dailynews.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent — runs the daily news agent every morning at 9:00 AM.
  Installed to: ~/Library/LaunchAgents/com.achuthkartha.dailynews.plist
  
  To manually load:   launchctl load ~/Library/LaunchAgents/com.achuthkartha.dailynews.plist
  To manually unload: launchctl unload ~/Library/LaunchAgents/com.achuthkartha.dailynews.plist
  To run right now:   launchctl start com.achuthkartha.dailynews
-->
<plist version="1.0">
<dict>

  <!-- Unique identifier for this job -->
  <key>Label</key>
  <string>com.user.dailynews</string>

  <!-- The command to run -->
  <key>ProgramArguments</key>
  <array>
    <string>PYTHON_PATH_PLACEHOLDER</string>
    <string>SCRIPT_PATH_PLACEHOLDER</string>
    <string>--quiet</string>
  </array>

  <!-- Run at 9:00 AM every day -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Environment variables (Ollama must be reachable) -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>HOME_PATH_PLACEHOLDER</string>
  </dict>

  <!-- Log output here if something goes wrong -->
  <key>StandardOutPath</key>
  <string>LOG_PATH_PLACEHOLDER/daily-news-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>LOG_PATH_PLACEHOLDER/daily-news-stderr.log</string>

  <!-- Run even if we missed the window (e.g. Mac was asleep) -->
  <key>RunAtLoad</key>
  <false/>

</dict>
</plist>
